# endpoints/static-visualization-server/storeconnect-owldoc.service
[Unit]
Description=StoreConnect's static visualization server
After=docker.service

[Service]
Type=simple
RemainAfterExit=yes
Restart=on-failure

# Server-specific configuration
User=administrateur
ExecStart=/usr/local/bin/docker-compose up
WorkingDirectory=/opt/storeconnect-ontology/endpoints/static-visualization-server

[Install]
WantedBy=multi-user.target